# templates/etc/init/openconnect.conf.erb
description "openconnect vpn client daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 120

<% if @dnsupdate == false -%>
env DNS_UPDATE=no
export DNS_UPDATE
<% end -%>

exec cat /etc/openconnect/network.passwd | openconnect \
  <%= @url -%> \
  --syslog \
  --no-dtls \
  --script /etc/vpnc/vpnc-script \
<% unless @cacerts.empty? -%>
  --CAfile /etc/openconnect/network.cacerts \
<% end -%>
<% if @authgroup -%>
  --authgroup=<%= @authgroup %>
<% end -%>
<% unless @servercert.empty? -%>
  --servercert <%= @servercert -%> \
<% end -%>
  --user <%= @user -%> \
  --passwd-on-stdin
